# tasks.crontab
* * * * * cd /usr/src/app && npm run exec unknown-task